# server/launchd/com.personalai.backup-nightly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.personalai.backup-nightly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/aspirebunny/Library/Application Support/personal-ai-assistant/scripts/pai_backup_nightly.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/aspirebunny/Desktop/personal-ai-assistant</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PAI_PROJECT_ROOT</key>
    <string>/Users/aspirebunny/Desktop/personal-ai-assistant</string>
    <key>PAI_SERVER_DIR</key>
    <string>/Users/aspirebunny/Desktop/personal-ai-assistant/server</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/pai_backup.launchd.out</string>

  <key>StandardErrorPath</key>
  <string>/tmp/pai_backup.launchd.err</string>
</dict>
</plist>
